# templates/btrfs_snapshot.timer
[Unit]
Description=nightly snapshot

[Timer]
OnCalendar=*-*-* 00:00:00
Unit=btrfs_snapshot.service

[Install]
WantedBy=timers.target
