# roles/os_customisation/templates/btrfs_snapshot@.service
[Unit]
Description=Snapshot BTRFS

[Service]
Type=oneshot
ExecStart=/root/code/interval_btrfs_snapshot.sh %I

[Install]
WantedBy=multi-user.target
