# roles/os_customisation/templates/btrfs_snapshot.service
[Unit]
Description=Snapshot BTRFS

[Service]
Type=oneshot
ExecStart=/root/code/daily_btrfs_snapshot.sh

[Install]
WantedBy=multi-user.target
